# wakelock/wakelock.service
[Unit]
Description=Wakelock

[Service]
Type=simple
ExecStart=/usr/local/bin/wakelock inhibit
