# wakelock/wakelock-detect.service
[Unit]
Description=Wakelock boot-time hook
After=dbus.service

[Service]
ExecStart=/usr/local/bin/wakelock

[Install]
WantedBy=multi-user.target
